# sendy.cron
*/5 * * * * curl -s "http://localhost/scheduled.php?cron_key=${CRON_KEY}" > /dev/null 2>&1
